# cnb-builder-shim/rootfs/dev_sandbox/etc/services.d/code-server/run
#!/usr/bin/with-contenv bash

echo "starting code-server"

/usr/bin/code-server --bind-addr 0.0.0.0:8080 --disable-telemetry --disable-update-check
