# cnb-builder-shim/rootfs/devsandbox/etc/s6-overlay/s6-rc.d/dev-server/run
#!/usr/bin/with-contenv bash

echo "starting dev-builder"

/cnb/devsandbox/bin/dev-entrypoint
